# docker/rootfs/etc/systemd/system/user@.service.d/resources.conf
[Service]
CPUWeight=100
OOMScoreAdjust=0
